# deploy/fingerprint-logger.service
[Unit]
Description=Fingerprint Logger - Pi5 Brain
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/esp32-autopilot-driver-monitor
Environment=APP_CONFIG=/home/pi/esp32-autopilot-driver-monitor/config/config.toml
ExecStart=/home/pi/esp32-autopilot-driver-monitor/.venv/bin/python /home/pi/esp32-autopilot-driver-monitor/pc/launch_fingerprint_logger.py
Restart=always
RestartSec=2
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
